# max-mate-settings/etc/xdg/autostart/xdg-user-dirs-mate.desktop
[Desktop Entry]
Version=1.0
Encoding=UTF-8
Name=XDG User Directories in MATE
TryExec=xdg-user-dirs-update
Exec=xdg-user-dirs-update --force
Icon=folder
Terminal=false
Comment=Create XDG user directories in MATE
Type=Application
OnlyShowIn=MATE;
