# max-multiseat/lib/udev/rules.d/50-usbseat.rules
########################################### NEW #################################################
# New MAX 6.0 rules use /tmp/seat.db to save BUSNUM/DEVNUM with a SEAT_ID
# See at /var/log/syslog | grep make-usbseat for debug output

# old HP multiseat devices
KERNEL=="fb*",SUBSYSTEMS=="usb",PROGRAM="/lib/udev/make-usbseat.sh %p/../../ displayfb",\
                                SYMLINK+="usbseat/%c/display",RUN+="usbseat.sh %c displayfb"

# MWS300
ATTR{idVendor}=="0711",ATTR{idProduct}=="5100",PROGRAM="/lib/udev/make-usbseat.sh %p display",\
                                SYMLINK+="usbseat/%c/display",ENV{SEAT_ID}="%c",RUN+="usbseat.sh %c display"

# MWS300v2
ATTR{idVendor}=="0711",ATTR{idProduct}=="5500",PROGRAM="/lib/udev/make-usbseat.sh %p display",\
                                SYMLINK+="usbseat/%c/display",ENV{SEAT_ID}="%c",RUN+="usbseat.sh %c display"



KERNEL=="mouse*", SUBSYSTEMS=="usb", ATTRS{bInterfaceClass}=="03", ATTRS{bInterfaceProtocol}=="02", \
        PROGRAM="/lib/udev/make-usbseat.sh %p/../../../../ mouse",SYMLINK+="usbseat/%c/mouse",RUN+="usbseat.sh %c mouse"

KERNEL=="event*", SUBSYSTEM=="input", ATTRS{bInterfaceClass}=="03", ATTRS{bInterfaceProtocol}=="01",\
        PROGRAM="/lib/udev/make-usbseat.sh %p/../../../../ keyboard",SYMLINK+="usbseat/%c/keyboard",RUN+="usbseat.sh %c keyboard"

KERNEL=="control*", SUBSYSTEM=="sound", SUBSYSTEMS=="usb", \
        PROGRAM="/lib/udev/make-usbseat.sh %p/../../../../ sound",SYMLINK+="usbseat/%c/sound"


#################################################################################################
