# max-mate-settings/usr/share/mate/autostart/max-mate-settings-overlay.desktop
[Desktop Entry]
Type=Application
Exec=/usr/lib/max-mate/max-mate-settings-overlay
Terminal=false
X-MATE-Autostart-enabled=true
Name[en_GB]=Settings Overlay
Name=Settings Overlay
Comment[en_GB]=MAX MATE Settings Overlay
Comment=MAX MATE Settings Overlay
