# max-mate-settings/usr/share/mate/autostart/plank.desktop
[Desktop Entry]
Type=Application
Exec=/usr/bin/plank
Terminal=false
Name=Top launcher for infantil
Name[es]=Barra superior para infantil
NoDisplay=true
X-MATE-Autostart-enabled=true
X-MATE-Autostart-Delay=4
